# cron/prod.cron
0 0 * * * php /home/nikita/projects/yii2/cron/prod.cron